# docker/frankenphp/s6/frankenphp/run
#!/command/with-contenv sh
# Load the secrets generated by init.sh (.env.local) into the process
# environment — Caddy reads {env.MERCURE_JWT_SECRET} at startup.
if [ -f /var/www/html/.env.local ]; then
  set -a
  . /var/www/html/.env.local
  set +a
fi
# Caddy needs a writable HOME for its internal data store (XDG paths).
# var/caddy is created and chowned by init.sh.
export HOME=/var/www/html/var/caddy

# ── FrankenPHP worker mode (opt-in, default OFF) ────────────────────────
# Worker mode boots the Symfony kernel once and keeps it resident, so requests
# skip the per-request kernel bootstrap. Controlled by PRISMARR_WORKER:
#   unset / "" / 0 / false / no / off  → OFF (classic per-request mode)
#   1  (or true / yes / on)            → ON, FrankenPHP's default count (2x CPU)
#   a positive integer N (e.g. 4)      → ON, with N resident workers
# PRISMARR_WORKER_NUM, if set, overrides the count in every case. The directive
# is injected into the global `frankenphp { }` block via the
# {$PRISMARR_WORKER_DIRECTIVE} Caddyfile placeholder; an unset/OFF value leaves
# it empty — byte-for-byte the classic path. Symfony 8's default runtime
# auto-detects FrankenPHP's worker loop (FRANKENPHP_WORKER) — no APP_RUNTIME /
# extra package needed.
case "${PRISMARR_WORKER:-}" in
  ''|0|false|False|FALSE|no|No|NO|off|Off|OFF)
    : ;;  # classic — leave the directive empty
  *)
    _pw_num="${PRISMARR_WORKER_NUM:-}"
    # No explicit count and PRISMARR_WORKER is a plain integer >1 → use it as
    # the worker count (so PRISMARR_WORKER=4 means "4 workers").
    case "${PRISMARR_WORKER}" in
      ''|*[!0-9]*) : ;;                                                       # not a pure integer
      *) [ -z "${_pw_num}" ] && [ "${PRISMARR_WORKER}" -gt 1 ] && _pw_num="${PRISMARR_WORKER}" ;;
    esac
    if [ -n "${_pw_num}" ]; then
      export PRISMARR_WORKER_DIRECTIVE="worker /var/www/html/public/index.php ${_pw_num}"
    else
      export PRISMARR_WORKER_DIRECTIVE="worker /var/www/html/public/index.php"
    fi
    ;;
esac

# Drop privileges: FrankenPHP and all PHP scripts run as www-data (UID 82).
exec s6-setuidgid www-data frankenphp run --config /etc/caddy/Caddyfile
