# docker/frankenphp/s6/frankenphp/run
#!/command/with-contenv sh
# Load the secrets generated by init.sh (.env.local) into the process
# environment — Caddy reads {env.MERCURE_JWT_SECRET} at startup.
if [ -f /var/www/html/.env.local ]; then
  set -a
  . /var/www/html/.env.local
  set +a
fi
# Caddy needs a writable HOME for its internal data store (XDG paths).
# var/caddy is created and chowned by init.sh.
export HOME=/var/www/html/var/caddy

# ── FrankenPHP worker mode (opt-in, default OFF) ────────────────────────
# PRISMARR_WORKER=1 turns on Symfony worker mode: the kernel is booted once
# and kept resident, so requests skip the per-request kernel bootstrap.
# PRISMARR_WORKER_NUM optionally pins the worker/thread count (default is
# FrankenPHP's own 2x CPU). The directive is injected into the global
# `frankenphp { }` block via the {$PRISMARR_WORKER_DIRECTIVE} placeholder in
# the Caddyfile; leaving PRISMARR_WORKER unset keeps the classic per-request
# path untouched. Symfony's default runtime auto-detects FrankenPHP's worker
# loop (FRANKENPHP_WORKER) — no APP_RUNTIME / extra package needed on SF 8.
if [ "${PRISMARR_WORKER:-}" = "1" ]; then
  if [ -n "${PRISMARR_WORKER_NUM:-}" ]; then
    export PRISMARR_WORKER_DIRECTIVE="worker /var/www/html/public/index.php ${PRISMARR_WORKER_NUM}"
  else
    export PRISMARR_WORKER_DIRECTIVE="worker /var/www/html/public/index.php"
  fi
fi

# Drop privileges: FrankenPHP and all PHP scripts run as www-data (UID 82).
exec s6-setuidgid www-data frankenphp run --config /etc/caddy/Caddyfile
